# nodes/proxmox/sanoid/syncoid-backup.timer
[Unit]
Description=Daily syncoid replication

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
